# external/linux/sndio-1.5.0/contrib/sndiod.service.in
[Unit]
Description=sndio audio and MIDI server
After=network.target

[Service]
Type=forking
Restart=on-abort
EnvironmentFile=-/etc/default/sndiod
ExecStart=@bindir@/sndiod $DAEMON_OPTS

[Install]
WantedBy=multi-user.target
